# deploy/systemd/factory-telemetry.service
[Unit]
Description=Agentryx Dev Factory - Telemetry + cognitive-engine spawner (port 4401)
After=network.target

[Service]
Type=simple
User=subhash.thakur.india
Group=subhash.thakur.india
WorkingDirectory=/home/subhash.thakur.india/Projects/agentryx-factory/factory-dashboard
Environment=PATH=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Phase 2D — activate the LLM router for spawned cognitive-engine children.
# Default backend is openrouter; flip to litellm by changing LLM_ROUTER_BACKEND.
# Provider keys are sourced from claw-code-parity/.env (already chmod 600).
Environment=USE_ROUTER=true
Environment=LLM_ROUTER_BACKEND=openrouter
EnvironmentFile=-/home/subhash.thakur.india/Projects/claw-code-parity/.env

ExecStart=/home/subhash.thakur.india/.nvm/versions/node/v22.22.2/bin/node server/telemetry.mjs
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
